# root/etc/periodic/15min/15min.sh
#!/bin/bash
echo "`date`: /etc/periodic/15min/15min.sh" >> /tmp/15min.log
